# auto_start/SensorCleanUpgradeSMB.service
[Unit]
Description=KootNet Sensors - Clean Upgrade Process SMB

[Service]
Type=simple
Restart=no
User=root
Group=root
ExecStart=/bin/bash /home/kootnet_data/scripts/clean_upgrade_smb.sh

[Install]
WantedBy=multi-user.target
