# extras/deb_package_files/etc/systemd/system/SensorUpgradeOnlineCleanDEV.service
[Unit]
Description=KootNet Sensors - Clean HTTP Upgrade

[Service]
Type=simple
Restart=no
User=root
Group=root
ExecStart=/bin/bash /home/kootnet_data/scripts/clean_upgrade_http.sh dev

[Install]
WantedBy=multi-user.target
