# extras/deb_package_files/DEBIAN/prerm
#!/bin/bash
systemctl disable KootnetSensors
systemctl stop KootnetSensors
crontab -u root -r
find /opt/kootnet-sensors -type d -name __pycache__ -exec rm -rf {} \; 2>/dev/null
exit 0
