# extras/deb_package_files/DEBIAN/prerm
#!/bin/bash
systemctl disable KootnetSensors
systemctl stop KootnetSensors
crontab -u root -r
exit 0
